# ops/oracle/postgres-backup.timer
[Unit]
Description=Daily Gestao Facil Postgres Backup

[Timer]
OnCalendar=*-*-* 03:35:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
